# install/systemd/magicka.service
[Unit]
Description=Posda Lanterna Magicka
After=network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/oneposda/install/start-magicka.sh

[Install]
WantedBy=multi-user.target
