# install/systemd/posda-ffp.service
[Unit]
Description=Posda FastFileProcessor
After=network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/oneposda/install/start-ffp.sh

[Install]
WantedBy=multi-user.target
