# web/quince/server/quince-server.service
[Unit]
Description=Quince Server
After=network.target posda.service

[Service]
User=posda
Group=posda
WorkingDirectory=/home/posda/quince/server
ExecStart=/opt/python36/bin/python3.6 /home/posda/quince/server/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
