# posda/supervisord_scripts/posda-file-process.ini
[program:posda-file-process]
command=/home/posda/posdatools/systemd/run_in_posda_env.sh systemd/start_PosdaFileProcessDaemon.pl
directory=/home/posda/posdatools/systemd
autostart=true

stdout_events_enabled = true
stderr_events_enabled = true
