# posda/supervisord_scripts/posda.ini
[program:posda]
command=/home/posda/posdatools/posda_start.sh
directory=/home/posda/posdatools
autostart=true

stdout_events_enabled = true
stderr_events_enabled = true
